# crontab
#!/bin/sh

/usr/bin/curl https://cogether.azurewebsites.net/help-center/freq/ >> /proc/$(cat /var/run/crond.pid)/fd/1 2>&1

