# crontab
1 * * * * /usr/bin/curl -X -v GET https://www.google.com/
#